# crontab
0 5 * * * /root/backup.sh > /dev/null 2>&1
